# config/sera/scripts/sera-esp32-bridge.user.service
[Unit]
Description=OpenClaw -> ESP32-C6 LCD status bridge (Sera) [user]

[Service]
Type=simple
WorkingDirectory=/home/pi/.openclaw/workspace
ExecStart=/usr/bin/env node /home/pi/.openclaw/workspace/scripts/esp32-status-bridge.mjs
Restart=always
RestartSec=2
Environment=SERA_PROJECT=Sera
Environment=OPENCLAW_LOG_DIR=/tmp/openclaw

[Install]
WantedBy=default.target
