# config/openclaw/scripts/vibemon-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibemon.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec /usr/bin/env node HOME_DIR/.openclaw/workspace/scripts/vibemon-bridge.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>HOME_DIR/.openclaw/workspace</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PROJECT_NAME</key>
        <string>OpenClaw</string>
        <key>OPENCLAW_LOG_DIR</key>
        <string>/tmp/openclaw</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>HOME_DIR/.openclaw/logs/vibemon-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>HOME_DIR/.openclaw/logs/vibemon-bridge.error.log</string>
</dict>
</plist>
