# config/openclaw/scripts/vibemon-bridge.user.service
[Unit]
Description=OpenClaw -> ESP32-C6 LCD status bridge (OpenClaw) [user]

[Service]
Type=simple
WorkingDirectory=/home/pi/.openclaw/workspace
ExecStart=/usr/bin/env node /home/pi/.openclaw/workspace/scripts/vibemon-bridge.mjs
Restart=always
RestartSec=2
Environment=PROJECT_NAME=OpenClaw
Environment=OPENCLAW_LOG_DIR=/tmp/openclaw

[Install]
WantedBy=default.target
